# docker/backend/files/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=/usr/local/bin/node /server.js
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
